# emby.conf
[program:emby]
user = nobody
command = /home/nobody/setup.sh