# emby.conf
[program:emby]
user = nobody
command = /usr/bin/emby-server